# config/supervisord.conf
[program:myproject]
command = /home/ubuntu/myproject/bin/gunicorn_start
user = ubuntu
stdout_logfile = /home/ubuntu/myproject/logs/supervisor.log
redirect_stderr = true
